# templates/etcd-install.service
[Unit]
Requires={{etcd_service}}
After={{etcd_service}}

[Service]
ExecStart=/usr/bin/sleep infinity

[Install]
WantedBy=multi-user.target
